# supervisor-latest-tag/supervisord-sickrage.ini
[program:sickrage]
command=/usr/bin/env python2 /sickrage/SickBeard.py --datadir=/config --nolaunch
user=sickrage
autostart=true
autorestart=true
priority=1
umask=002
startretries=3
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
